# roles/snapshot_service/templates/create-snapshot.target.j2
# {{ ansible_managed }}
[Unit]
Description=Create LVM snapshots
Requires=multi-user.target
