# roles/snapshot_service/templates/restore-snapshot.target.j2
# {{ ansible_managed }}
[Unit]
Description=Restore LVM snapshots
Requires=multi-user.target
