# roles/snapshot_service/templates/restore-snapshot.target.j2
[Unit]
Description=Restore LVM snapshots
Requires=multi-user.target
